# target/linux/rockchip/armv8/base-files/etc/hotplug.d/usb/15-usb-wifi-config
#!/bin/sh

case "$(board_name)" in
rockchip,rk3568-photonicat)
    exit
    ;;
esac

if [ "${ACTION}" = "add" ]; then
  usbmode -s && sleep 5
  [ -n "$(cat /proc/net/wireless | grep wlan)" ] && wifi config && wifi
fi
